# config/flasktodo.service
[Unit]
Description=uWSGI Container Server for Flask Todo
After=network.target

[Service]
User=ec2-user
Group=nginx
WorkingDirectory=/var/www
Environment="PATH=/var/www/venv/bin"
ExecStart=/var/www/venv/bin/uwsgi --ini config/flasktodo.ini

[Install]
WantedBy=multi-user.target